# configs/supervisor.ini
[program:cloudping]
environment=NEW_RELIC_CONFIG_FILE=/home/ubuntu/cloudping/configs/newrelic.ini
directory=/home/ubuntu/cloudping/configs
command=newrelic-admin run-program /usr/local/bin/uwsgi --ini uwsgi.ini
autostart=true
autorestart=true
stopsignal=QUIT
killasgroup=true
stdout_logfile=/home/ubuntu/cloudping/logs/uwsgi.stdout.log
stderr_logfile=/home/ubuntu/cloudping/logs/uwsgi.stderr.log
